# freeipa-letsencrypt.service
[Unit]
Description=FreeIPA Let's Encrypt service
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/freeipa-letsencrypt-apply-certs
