# freeipa-letsencrypt.path
[Path]
PathModified=/etc/letsencrypt/live/%H

[Install]
WantedBy=multi-user.target
